# server/compendia-server.service
[Unit]
Description = Runs the Compendia server
[Service]
ExecStart = /usr/share/Compendia/server/start.sh
[Install]
WantedBy = multi-user.target

